# src/yolov3_parking/src/docker-novnc/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 8080 localhost:5900
autorestart=true
