# src/yolov3_parking/src/docker-novnc/conf.d/xterm.conf
[program:xterm]
command=xterm
autorestart=true
